# crons/orchestrator/backup_daily_brt.cron
# Cron schedule (BRT): daily at 00:00
CRON_TZ=America/Sao_Paulo
0 0 * * * /local/crons/orchestrator/backup_daily_brt.sh
